# buildout_templates/supervisor-django-celery.conf.in
[program:${name}.django-celery]
process_name=%(program_name)s_%(process_num)s
command=${buildout:bin-directory}/${control-script} celery worker -B --schedule=/tmp/${name}-celerybeat-schedule --concurrency=1
user=${user}
stdout_logfile=${log}
stderr_logfile=${error-log}
stopsignal=QUIT
environment=PYTHON_EGG_CACHE='/var/www/.python-eggs'
