# buildout_templates/supervisor-haproxy.conf.in
[program:${name}.haproxy]
command=haproxy -f ${config} -sf
stdout_logfile=${log}
redirect_stderr=true
stopsignal=QUIT
