# buildout_templates/supervisor-django.conf.in
[program:${name}.django]
process_name=%(program_name)s
command=${buildout:bin-directory}/${control-script} run_gunicorn -w ${num-processes} --max-requests=${max-requests} -b ${host}:${port-prefix}1
user=${user}
stdout_logfile=${log}
stderr_logfile=${error-log}
stopsignal=QUIT
environment=PYTHON_EGG_CACHE='/var/www/.python-eggs'
