# buildout_templates/supervisor-haproxy.conf.in
[program:${name}.haproxy]
command=haproxy -f ${config} -sf
stdout_logfile=${log}
stdout_errfile=${error-log}
stopsignal=QUIT
