# conf/supervisor.conf.in
[program:skeleton-admin-django]
process_name=%(program_name)s
command={{ deploy_root }}/python/bin/gunicorn wsgi -w 4 --max-requests=10000 --bind 127.0.0.1:9091
directory={{ deploy_root }}/{{ app_name }}/
user=www-data
stdout_logfile={{ deploy_root }}/log/{{ app_name }}-django-admin.log
stderr_logfile={{ deploy_root }}/log/{{ app_name }}-django-admin-error.log
stopsignal=QUIT
environment=
    PYTHON_EGG_CACHE='/var/www/.python-eggs',
    DJANGO_SETTINGS_MODULE=project.settings_live_admin_site

[program:skeleton-basic-django]
process_name=%(program_name)s
command={{ deploy_root }}/python/bin/gunicorn wsgi -w 4 --max-requests=10000 --bind 127.0.0.1:9021
directory={{ deploy_root }}/{{ app_name }}/
user=www-data
stdout_logfile={{ deploy_root }}/log/{{ app_name }}-django-basic.log
stderr_logfile={{ deploy_root }}/log/{{ app_name }}-django-basic-error.log
stopsignal=QUIT
environment=
    PYTHON_EGG_CACHE='/var/www/.python-eggs',
    DJANGO_SETTINGS_MODULE=project.settings_live_basic_site

[program:skeleton-smart-django]
process_name=%(program_name)s
command={{ deploy_root }}/python/bin/gunicorn wsgi -w 4 --max-requests=10000 --bind 127.0.0.1:9031
directory={{ deploy_root }}/{{ app_name }}/
user=www-data
stdout_logfile={{ deploy_root }}/log/{{ app_name }}-django-smart.log
stderr_logfile={{ deploy_root }}/log/{{ app_name }}-django-smart-error.log
stopsignal=QUIT
environment=
    PYTHON_EGG_CACHE='/var/www/.python-eggs',
    DJANGO_SETTINGS_MODULE=project.settings_live_smart_site

[program:skeleton-web-django]
process_name=%(program_name)s
command={{ deploy_root }}/python/bin/gunicorn wsgi -w 4 --max-requests=10000 --bind 127.0.0.1:9011
directory={{ deploy_root }}/{{ app_name }}/
user=www-data
stdout_logfile={{ deploy_root }}/log/{{ app_name }}-django-web.log
stderr_logfile={{ deploy_root }}/log/{{ app_name }}-django-web-error.log
stopsignal=QUIT
environment=
    PYTHON_EGG_CACHE='/var/www/.python-eggs',
    DJANGO_SETTINGS_MODULE=project.settings_live_web_site

[program:skeleton-deviceproxy]
command={{ deploy_root }}/python/bin/twistd --pidfile=/tmp/deviceproxy.pid -n devproxy --config {{ deploy_root }}/{{ app_name }}/deviceproxy_site.yaml
directory={{ deploy_root }}/{{ app_name }}/
user=www-data
numprocs=1
stdout_logfile={{ deploy_root }}/log/{{ app_name }}-deviceproxy.log
stderr_logfile={{ deploy_root }}/log/{{ app_name }}-deviceproxy-error.log
stopsignal=QUIT
startsecs=10
environment=PYTHON_EGG_CACHE='/var/www/.python-eggs'

[program:skeleton-celery-worker]
process_name=%(program_name)s_%(process_num)s
command={{ deploy_root }}/{{ app_name }} manage.py celery worker
directory={{ deploy_root }}/{{ app_name }}/
user=www-data
stdout_logfile={{ deploy_root }}/log/{{ app_name }}-django-celery.log
stderr_logfile={{ deploy_root }}/log/{{ app_name }}-django-celery-error.log
stopsignal=QUIT
environment=PYTHON_EGG_CACHE='/var/www/.python-eggs'

[program:skeleton-celery-beat]
process_name=%(program_name)s_%(process_num)s
command={{ deploy_root }}/{{ app_name }} manage.py celery beat --pidfile=/tmp/{{ app_name }}-django-celery-beat.pid --schedule=/tmp/{{ app_name }}-celery-beat-schedule
directory={{ deploy_root }}/{{ app_name }}/
user=www-data
stdout_logfile={{ deploy_root }}/log/{{ app_name }}-django-celery.log
stderr_logfile={{ deploy_root }}/log/{{ app_name }}-django-celery-error.log
stopsignal=QUIT
environment=PYTHON_EGG_CACHE='/var/www/.python-eggs'

[group:skeleton]
programs=skeleton-admin-django,skeleton-basic-django,skeleton-smart-django,skeleton-web-django,skeleton-deviceproxy,skeleton-celery-worker,skeleton-celery-beat
