# scripts/resources/supervisor.haproxy.conf
[program:haproxy]
command=/usr/local/bin/haproxy-start-wrapper.sh
user=www-data
stdout_logfile=/var/praekelt/log/haproxy.log
redirect_stderr=true
stopsignal=QUIT
