# scripts/resources/supervisor.haproxy.conf
[program:haproxy]
command=/usr/local/bin/haproxy-start-wrapper.sh
user=www-data
stopsignal=QUIT
